# supervisord.conf
; .....

;go-captcha-service
[program:go-captcha-service]
command=/app/go-captcha-service-[xxx]
; ...

; .....